# server/crash/agora-crash.service
[Unit]
Description=Agora anonymous crash-report receiver
After=network.target

[Service]
Type=simple
User=newoether
Group=newoether
ExecStart=/usr/bin/python3 /opt/agora-crash/agora-crash.py
Restart=on-failure
RestartSec=2

# Sandboxing — the service only needs to bind loopback and write its data dir.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/agora-crash
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
